# Caeli-server.service
[Unit]
Description=Serveur Caeli

[Service]
ExecStart=/usr/local/bin/Caeli_Serveur

[Install]
WantedBy=Caeli-client.timer
